# work/unsquashing/squashfs-root/etc/inittab
::wait:/etc/rc.d/rc
::sysinit:/etc/initsh
::shutdown:/sbin/shutdown h
::respawn:/sbin/getty 57600 ttyS1
