# work/unsquashing/squashfs-root/etc/udev/rules.d/10-net.rules
SUBSYSTEM=="gmac_link_state", ACTION=="online", RUN+="/etc/init.d/network &"
