# work/unsquashing/squashfs-root/etc/udev/rules.d/15-usb_vbus.rules
SUBSYSTEM=="reset_button",  ACTION=="add", RUN+="/etc/udev/script/reset.sh"
SUBSYSTEM=="power_button",  ACTION=="add", RUN+="/etc/udev/script/power.sh"
